# config/launchd/com.local.jarvis.workbench.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.local.jarvis.workbench</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/YOUR_USER/Jarvis/services/start-workbench.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USER/Jarvis/logs/workbench.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USER/Jarvis/logs/workbench.err.log</string>
  </dict>
</plist>
